# systemd/eg-buswatch.service
[Unit]
Description=Evergreeen Bus Watch Service

[Service]
Type=simple
User=opensrf
Group=opensrf
ExecStart=/usr/local/bin/eg-buswatch
