# evergreen-bin/systemd/eg-websockets.service
[Unit]
Description=Evergreen Websockets

[Service]
Type=simple
User=opensrf
Group=opensrf
ExecStart=/usr/local/bin/eg-websockets
